# deploy/sms-viewer.service
[Unit]
Description=NexSMS SMS Viewer
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/app
EnvironmentFile=/root/app/.env
ExecStart=/root/.nvm/versions/node/v24.14.1/bin/node /root/app/server.mjs
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
